# bumble-bee/root.crontab
0 3 * * * bash /home/lbouch/infastructure/bumble-bee/scripts/backup-notes.sh
0 3 * * * bash /home/lbouch/infastructure/bumble-bee/scripts/backup-financial.sh
0 3 * * * bash /home/lbouch/infastructure/bumble-bee/apps/timetagger/backup.sh
0 3 * * * bash /home/lbouch/infastructure/bumble-bee/apps/navidrome/backup.sh
0 4 * * * bash /home/lbouch/infastructure/bumble-bee/apps/immich/backup.sh
0 3 * * * bash /home/lbouch/infastructure/bumble-bee/apps/actual-budget/backup.sh
0 3 * * * bash /home/lbouch/infastructure/bumble-bee/apps/grocy/backup-grocy.sh
